# specs/icecast-kh/SOURCES/icecast.service
[Unit]
Description=Icecast Network Audio Streaming Server
After=network.target

[Service]
ExecStart=/usr/bin/icecast -c /etc/icecast.xml
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
