# vuls/SOURCES/goval-dictionary-fetch.cron
5 * * * * goval /usr/bin/goval-dictionary -dbpath="/opt/vuls/data" fetch-redhat 6 7 &> /dev/null || :
